# bob/etc/supervisor/supervisord.conf
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = /dev/shm/supervisor.sock
chmod = 0777
chown= nobody:nogroup

[supervisorctl]
serverurl = unix:///dev/shm/supervisor.sock
prompt = mysupervisor

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/run/supervisord.pid
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:newrelic-daemon]
command=/usr/bin/newrelic-daemon  --pidfile /var/run/newrelic-daemon.pid --logfile=/dev/stdout --foreground
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize -c /etc/php/7.0/fpm
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g "daemon off;"
priority=3
numprocs=1
user=root
directory=/app/
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx-php-fpm-prometheus]
command=/bin/nginx-php-fpm-prometheus -fpm-status-url=http://localhost/fpm_status_bob -nginx-status-url=http://localhost/nginx_status
priority=4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
